# extras/modelrelay-server.service
# ModelRelay — modelrelay-server systemd unit
#
# Install:
#   sudo cp modelrelay-server.service /etc/systemd/system/
#   sudo cp proxy.env.example /etc/modelrelay/proxy.env  # edit values
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now modelrelay-server

[Unit]
Description=ModelRelay central proxy server
Documentation=https://modelrelay.io
After=network.target

[Service]
Type=simple

# Run as a dedicated unprivileged user. Create it with:
#   sudo useradd --system --no-create-home modelrelay
User=modelrelay
Group=modelrelay

WorkingDirectory=/var/lib/modelrelay

# Load environment from file. The leading "-" means the service still starts
# if the file is missing (CLI defaults or env vars set elsewhere will apply).
EnvironmentFile=-/etc/modelrelay/proxy.env

ExecStart=/usr/local/bin/modelrelay-server

# Restart on crash but not on clean exit (e.g. SIGTERM during drain).
Restart=on-failure
RestartSec=5

# Hardening — drop capabilities the proxy doesn't need.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
